# examples/tg-refresh.cron
# Refresh the local Telegram cache every 15 minutes.
# Adjust the path to uv and the working directory for your environment.
*/15 * * * * cd /path/to/tg-cli && /path/to/uv run tg refresh --yaml >> /tmp/tg-refresh.log 2>&1
